# root/etc/services.d/nmbd/run
#!/usr/bin/env bash

set -eu
set -o pipefail


/usr/sbin/nmbd --foreground --no-process-group
