# root/etc/services.d/avahi/run
#!/usr/bin/env bash

set -eu
set -o pipefail


s6-svwait /var/run/s6/services/dbus


exec /usr/sbin/avahi-daemon --no-chroot
