# root/etc/services.d/sambd/run
#!/usr/bin/env bash

set -eu
set -o pipefail


exec /usr/sbin/smbd --foreground --no-process-group --debug-stdout
